# vendor/qcom/proprietary/fastmmi/res/config/le/start_mmi
#! /bin/sh
### BEGIN INIT INFO
# Provides:          start_ffbm
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: Start ffbm deamon
# Description:       This script enabls android debug bridge over usb.
#                    Does not start/stop on any runlevel. Called from upstart job only
### END INIT INFO
#
# Copyright (c) 2016 Qualcomm Technologies, Inc.
# All Rights Reserved.
# Confidential and Proprietary - Qualcomm Technologies, Inc.
#
# ffbm        init.d script to start ffbm daemon

set -e

case "$1" in
  start)
        echo -n "Starting mmi: "
        start-stop-daemon -S --exec /usr/bin/mmi -b start_mmi
        echo "done"
        ;;
  stop)
        echo -n "Stopping mmi: "
        start-stop-daemon -K -n mmi
        echo "done"
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  *)
        echo "Usage ffbm { start | stop | restart}" >&2
        exit 1
        ;;
esac

exit 0
