# vendor/qcom/proprietary/web-setup/avs-setup/mdns_avs_stop.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c '. /usr/share/mdns_avs_functions.sh && isAvsConfigured && systemctl stop mdns_avs || :'
