# poky/meta-qti-data-prop/recipes/data/files/qtid.service
[Unit]
SourcePath=/usr/bin/qti
Description= qti service

[Service]
Restart=always
ExecStart=/usr/bin/qti

[Install]
WantedBy=multi-user.target