# vendor/qcom/proprietary/drones/sensortest/config/start_adsp_sensor_test
#! /bin/sh

# Copyright (c) 2016 Qualcomm Technologies, Inc.
# All Rights Reserved.
# Confidential and Proprietary - Qualcomm Technologies, Inc.

### BEGIN INIT INFO
# Provides:          start adsp sensor test daemon
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: start adsp sensor test daemon
# Description:       This script start adsp sensor test daemon for FFBM.
#                    Does not start/stop on any runlevel. Called from upstart job only
### END INIT INFO

set -e

case "$1" in
  start)
        echo -n "Starting adsp sensor test daemon:"
        echo 1 > /sys/kernel/boot_slpi/boot
        start-stop-daemon -S -a /usr/bin/ffbm_adsp
        echo "done"
        ;;
  stop)
        echo -n "Stopping adsp sensor test daemon:"
        start-stop-daemon -K -n ffbm_adsp
        echo "done"
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  *)
        echo "Usage start_adsp_sensor_test { start | stop | restart}" >&2
        exit 1
        ;;
esac

exit 0
