# quectel-core/ql-manager/etc/ql_manager_server.service.in
[Unit]
Description=ql_manager_server
#Wants=netmgrd.service
After=qmuxd.service

[Service]
Type=simple
Restart=always
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/ql_manager_server
ExecStartPost=/bin/sh -c 'setprop sys.boot_completed 1'
RestartSec=2

[Install]
WantedBy=multi-user.target
