# data/netmgr/src/netmgrd.service.in
[Unit]
Description=Netmgrd Service

[Service]
Type=simple
ExecStart=/usr/bin/netmgrd
Restart=always

[Install]
WantedBy=multi-user.target
