# vendor/qcom/proprietary/fastmmi/qti-system-daemon/start_qti_system
#! /bin/sh
### BEGIN INIT INFO
# Provides:          start_ffbm
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: Start mmi deamon
# Description:       This script enabls mmi daemon.
#                    Does not start/stop on any runlevel. Called from upstart job only
### END INIT INFO
#
# Copyright (c) 2016 Qualcomm Technologies, Inc.
# All Rights Reserved.
# Confidential and Proprietary - Qualcomm Technologies, Inc.
#
# ffbm        init.d script to start ffbm daemon

set -e

case "$1" in
  start)
        echo -n "Starting QTI system daemon: "
        start-stop-daemon -S --exec /usr/bin/qti_system_daemon -b start_qti_system
        echo "done"
        ;;
  stop)
        echo -n "Stopping QTI system daemon: "
        start-stop-daemon -K -n qti_system
        echo "done"
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  *)
        echo "Usage ffbm { start | stop | restart}" >&2
        exit 1
        ;;
esac

exit 0
